# crons/sync-paperless.service
[Unit]
Description=Sync paperless documents to Google Drive

[Service]
ExecStart=%h/.local/bin/sync-paperless.sh
