# crons/sync-notes.service
[Unit]
Description=Sync notes to Google Drive

[Service]
ExecStart=%h/.local/bin/sync-notes.sh
